# flashable/data/rootfs-mount/etc/init/mount-android.conf
#
# Copyright (C) 2017-2020 The Halium Project
#
# Licensed under the Apache License, Version 2.0 (the "License");
# you may not use this file except in compliance with the License.
# You may obtain a copy of the License at
#
#      http://www.apache.org/licenses/LICENSE-2.0
#
# Unless required by applicable law or agreed to in writing, software
# distributed under the License is distributed on an "AS IS" BASIS,
# WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
# See the License for the specific language governing permissions and
# limitations under the License.
# Mount device partitions.
# Dualboot Script Mount.

description "Mount device partitions"

start on (startup
    and started udev)

task

script
    mkdir -p /dev/cpuset
    mount none /dev/cpuset -t cpuset -o nodev,noexec,nosuid

    while [ ! -e /dev/disk/by-partlabel/vendor ]
    do
        sleep 1
    done
       
#   mount /dev/disk/by-partlabel/system /android/system -t ext4 -o ro
#   mount /dev/disk/by-partlabel/vendor /vendor -t ext4 -o ro
    mount /userdata/vendor.img /vendor -t ext4 -o ro
    mount /dev/disk/by-partlabel/modem /vendor/firmware_mnt -t vfat -o ro,shortname=lower,uid=0,gid=1000,dmask=227,fmask=337
    mount /dev/disk/by-partlabel/bluetooth /vendor/bt_firmware -t vfat -o ro,shortname=lower,uid=1002,gid=3002,dmask=227,fmask=337
    mount /dev/disk/by-partlabel/dsp /vendor/dsp -t ext4 -o ro,nosuid,nodev,barrier=1
    mount /dev/disk/by-partlabel/persist /mnt/vendor/persist -t ext4 -o noatime,nosuid,nodev,barrier=1

    mount -o bind /android/data /data
    mount -o bind /android/system/system /system
    mount -o bind /android/system /var/lib/lxc/android/rootfs
end script
